# opt/local/etc/profile.d/mysql-profile.d.sh
#!/bin/sh
export PATH=$PATH:/opt/local/sbin/mysql/bin
export MANPATH=$MANPATH:/opt/local/sbin/mysql/man
