# opt/local/etc/profile.d/mysql-profile.d.sh
#!/bin/sh
export PATH=$PATH:/usr/local/sbin/mysql/bin
export MANPATH=$MANPATH:/usr/local/sbin/mysql/man
