# opt/local/etc/profile.d/php.sh
#!/bin/sh
export PATH=$PATH:/opt/local/sbin/php/bin