# opt/local/etc/profile.d/php5-profile.d.sh
#!/bin/sh
export PATH=$PATH:/usr/local/sbin/php/bin
